# deployment/supervisor/chef.conf
[program:chef]
command=/home/ubuntu/chef/venv/bin/gunicorn -b localhost:8000 -w 4 chef:app
directory=/home/ubuntu/chef
user=ubuntu
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
